# ports/sysutils/linux-crashplan/files/crashplan.in
#!/bin/sh
#
# $FreeBSD$

# PROVIDE: crashplan
# REQUIRE: DAEMON

. /etc/rc.subr

name="crashplan"
rcvar=${name}_enable

cp_root="%%PREFIX%%"
java_root="${cp_root}/share/crashplan/jre"
crashplan_bin="${cp_root}/share/crashplan/bin"
command="${crashplan_bin}/CrashPlanEngine"
pidfile="/var/run/CrashPlanEngine.pid"

crashplan_start() {
    /usr/bin/install -d /var/log/crashplan || exit 1

    sed -i '' "s|/bin/ps -eo 'pid,cmd'|/bin/ps -ewwo 'pid,command'|" ${command}
    sed -i '' 's|PIDFILE="$TARGETDIR/${NAME}.pid"|PIDFILE="/var/run/${NAME}.pid"|' ${command}
    ${cp_root}/bin/patchelf --print-rpath ${java_root}/bin/java | grep '$ORIGIN' > /dev/null
    if [ $? -eq 0 ]; then
        ${cp_root}/bin/patchelf --set-rpath ${java_root}/lib/i386/jli ${java_root}/bin/java
    fi
    /usr/bin/cpuset -l 0 ${command} start
}

start_cmd="${name}_start"
stop_cmd="${command} stop"
status_cmd="${command} status"

load_rc_config ${name}
run_rc_command "$1"
